# crontab
* * * * * /opt/bitnami/php/sbin/php-fpm /app/vendor/yiisoft queue/run
